# installer/systemd/mrssh-backup.service
[Unit]
Description=MRSSH Daily Backup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mrssh-backup
